# route53.conf
#
description     "register in route53"

start on (started networking or starting ssh)

task
exec /usr/local/bin/update-zone.sh
